# templates/etc/systemd/system/docker.service.d/service-overrides-debian.conf.erb
[Service]
EnvironmentFile=-/etc/default/docker
EnvironmentFile=-/etc/default/docker-storage
ExecStart=
ExecStart=/usr/bin/<% if @daemon_command %><%= @daemon_command %><% else %><%= @docker_command %> <%= @daemon_subcommand %><% end -%> $OPTIONS \
        $DOCKER_STORAGE_OPTIONS
